# .claude/scripts/com.weekly-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.weekly-report</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- Update this path to match your local clone location -->
        <string>/Users/YOUR_USERNAME/path/to/weekly-report/.claude/scripts/run-weekly-report.sh</string>
    </array>
    <!-- Every Friday at 2pm — adjust Weekday (1=Mon…5=Fri) and Hour as needed -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>5</integer>
        <key>Hour</key>
        <integer>14</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/weekly-report-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/weekly-report-launchd.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
